# domain.service
[Unit]
Description=Domain Server Daemon

[Service]
ExecStart=domain serve

[Install]
WantedBy=default.target
